# build/File_Organizer.service
# Automatically generated by systemd-xdg-autostart-generator

[Unit]
Documentation=man:Organizador de arquivos da pasta de Download
SourcePath=/home/$USER/.config/autostart/File_Organizer.desktop
PartOf=graphical-session.target

Description= File_Organizer
After=graphical-session.target

[Service]
Type=exec
ExecStart=:/home/$USER/.local/share/file-organizer/dist/File_Organizer.desktop
Restart=no
TimeoutSec=5s
Slice=app.slice